# scripts/system-backup/systemd/system-backup-deep-check.timer
[Unit]
Description=Verificação profunda (read-data) semanal, domingos 04:00

[Timer]
OnCalendar=Sun *-*-* 04:00:00
RandomizedDelaySec=15m
Unit=system-backup-deep-check.service
Persistent=true

[Install]
WantedBy=timers.target
