# scripts/system-backup/systemd/system-backup.service
# Instalar: ver install-systemd.sh
[Unit]
Description=Backup do sistema (restic -> S3)
Documentation=file:/usr/local/lib/system-backup/README-RESTAURACAO.md
After=network-online.target time-sync.target
Wants=network-online.target

[Service]
Type=oneshot
# Variáveis e segredos (não fazer commit do ficheiro real)
EnvironmentFile=-/etc/system-backup/backup.env
ExecStart=/usr/local/lib/system-backup/01-backup.sh
# Só progride se o check + teste mínima de restauro passar
ExecStartPost=/usr/local/lib/system-backup/02-verify-restore.sh
Nice=10
IOSchedulingClass=idle
TimeoutStartSec=0
TimeoutStopSec=12h
StandardOutput=journal
StandardError=journal
SyslogIdentifier=system-backup

[Install]
WantedBy=multi-user.target
